# device/testing/device-kobo-clara/extract-waveform.initd
#!/sbin/openrc-run

description="Extract EPD waveform blob to /lib/firmware"

depend() {
	after root
}

start()
{
	# Only run once
	rm -f /etc/runlevels/*/$RC_SVCNAME

	dir="/lib/firmware/imx/epdc"

	ebegin "Copying EPD waveform to /lib/firmware..."
	mkdir -p "$dir" && cd "$dir" \
		&& extract-waveform.pl /dev/mmcblk0 epdc_PENG060D.fw \
		&& ln -sf epdc_PENG060D.fw epdc.fw
	eend $?
}
