# device/testing/device-samsung-p4wifi/90-device-samsung-p4wifi-audio.rules
ATTRS{id}=="wm8994aif1", RUN+="/usr/bin/alsaucm -c wm8994aif1 set _verb HiFi set _enadev Headset"