# device/testing/device-jolla-sbj/load-wlan-module
#!/sbin/openrc-run

description="Loading wlan module"

depend()
{
	need wcnss-wlan
	before wpa_supplicant
}

start()
{
	ebegin "Loading wlan module"
	modprobe wlan || return 1
}

stop()
{
	ebegin "Unloading wlan module"
	rmmod wlan || return 1
}

