# device/main/device-pine64-pinephone/60-flashlight.rules
SUBSYSTEM=="leds", DEVPATH=="*/*:flash", RUN+="/bin/chmod g+w /sys%p/brightness /sys%p/flash_strobe", RUN+="/bin/chown :video /sys%p/brightness /sys%p/flash_strobe"
