# device/main/device-pine64-pinephone/gpsd_pinephone.initd
#!/sbin/openrc-run

depend() {
	after gpsd
}

start() {
	gpsdctl add /dev/EG25.NMEA
}

stop() {
	gpsdctl remove /dev/EG25.NMEA
}
