# device/testing/bluetooth-raspberry-pi/pi-bluetooth.openrc
#!/sbin/openrc-run

command="/usr/bin/hciattach"
command_args="-n /dev/ttyAMA0 bcm43xx 115200 noflow -"
command_background=true
pidfile="/run/hciattach.pid"

depend() {
	before bluetooth
	after modules
	need localmount
}
