# main/hisi-wlan/hisi-wlan.initd
#!/sbin/openrc-run

description="Initialize hisi wlan driver"

depend()
{
	need udev
	before wpa_supplicant
}

start()
{
	ebegin "Enabling wlan driver"
	echo init > /sys/hisys/boot/plat
	echo init > /sys/hisys/boot/wifi
}
