# non-free/sgx-ddk-um/sgx-ddk-um.initd
#!/sbin/openrc-run

description="Initialize PowerVR SGX server"

depend() {
	need dev localmount
	after hwdrivers modules dev-mount
}

start() {
	for i in {0..5} ; do
		if /bin/pvrsrvinit ; then
			return 0
		fi
		sleep 5
	done
}

