# main/msm-fb-refresher/msm-fb-refresher.init
#!/sbin/openrc-run

command="/usr/sbin/msm-fb-refresher"
command_args="--loop"
supervisor="supervise-daemon"

depend() {
	after bootmisc
}
